# wsdd_0.6.4_all/DEBIAN/prerm
#!/bin/bash

systemctl --now disable wsdd
exit 0
